# etc/service/pdns/run
#!/bin/sh

exec 2>&1
/usr/sbin/pdns_server --config-dir=/home/user/config
